# conf/supervisord.conf
[program:taskserver]
command = /home/ken/.pyenv/shims/gunicorn sample:app --config /home/ken/apps/taskserver/conf/gunicorn.py
directory = /home/ken/apps/taskserver
user = ken
